# deploy/supervisord.conf
; =============================================================================
; LeAgent Supervisor Configuration
; =============================================================================
; Manages uvicorn workers with auto-restart and log rotation
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; LeAgent API Server (Gunicorn + Uvicorn)
; =============================================================================
[program:leagent-api]
command=gunicorn leagent.main:app
    --bind 0.0.0.0:8000
    --workers %(ENV_LEAGENT_WORKERS)s
    --worker-class uvicorn.workers.UvicornWorker
    --timeout 300
    --graceful-timeout 30
    --keep-alive 5
    --max-requests 10000
    --max-requests-jitter 1000
    --access-logfile -
    --error-logfile -
    --capture-output
    --enable-stdio-inheritance
directory=/app
user=leagent
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true

; Logging
stdout_logfile=/app/logs/leagent-api.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/app/logs/leagent-api-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
redirect_stderr=false

; Environment
environment=
    PYTHONUNBUFFERED="1",
    PYTHONFAULTHANDLER="1"

; =============================================================================
; Background Task Worker (Optional - for Celery/APScheduler)
; =============================================================================
[program:leagent-worker]
command=python -m leagent.cli.main worker
directory=/app
user=leagent
numprocs=1
autostart=false
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=60
stopsignal=TERM

stdout_logfile=/app/logs/leagent-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/leagent-worker-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

environment=
    PYTHONUNBUFFERED="1"

; =============================================================================
; Scheduler (APScheduler for cron jobs)
; =============================================================================
[program:leagent-scheduler]
command=python -m leagent.cli.main scheduler
directory=/app
user=leagent
numprocs=1
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM

stdout_logfile=/app/logs/leagent-scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/leagent-scheduler-error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; =============================================================================
; Process Groups
; =============================================================================
[group:leagent]
programs=leagent-api
priority=999

[group:leagent-background]
programs=leagent-worker,leagent-scheduler
priority=998

; =============================================================================
; Event Listeners (Optional - for crash notifications)
; =============================================================================
[eventlistener:process_exit]
command=python -c "import sys; sys.stdin.readline(); print('RESULT 2\\nOK')"
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=100
